# roles/gateway_common/files/etc/systemd/system/fastd@.service.d/override.conf
[Unit]
PartOf=fastd.service

[Service]
# Create these dirs before fastd startup
# for dignostics/metrics socket
ExecStartPre=mkdir -p /run/fastd
ExecStartPre=chown fastd:fastd /run/fastd
